# deployment/quickcommerce.service
[Unit]
Description=QuickCommerce FastAPI Application
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=notify
User=ubuntu
Group=ubuntu
WorkingDirectory=/var/www/quickcommerce
Environment="PATH=/var/www/quickcommerce/backend/venv/bin"
Environment="PYTHONPATH=/var/www/quickcommerce"
ExecStart=/var/www/quickcommerce/backend/venv/bin/uvicorn backend.main:app --host 127.0.0.1 --port 8080 --workers 2

# Restart configuration
Restart=always
RestartSec=3
StartLimitInterval=400
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quickcommerce

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
